# script/crontabs/root
* 4 * * *  sh /data/adb/modules/AD-Hosts/script/functions.sh